# scripts/launchd/com.myhouse.discover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 주간 신규편입 단지 탐색: 주 1회(월요일). 시간/요일 변경은 StartCalendarInterval 수정 후 재설치.
     Weekday: 0/7=일, 1=월 … 6=토. install_launchd.sh 가 __PLACEHOLDER__ 를 치환한다. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL_PREFIX__.discover</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>myhouse.cli</string>
    <string>discover</string>
    <string>--config</string>
    <string>__CONFIG__</string>
    <string>--trigger</string>
    <string>scheduled</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT__</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>__DISCOVER_HOUR__</integer>
    <key>Minute</key><integer>__DISCOVER_MINUTE__</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>Asia/Seoul</string>
    <key>PATH</key>
    <string>__VENV_BIN__:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/discover.out</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/discover.err</string>
</dict>
</plist>
