# scripts/launchd/com.myhouse.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 대시보드 상시 서버. install_launchd.sh 가 __PLACEHOLDER__ 를 치환한다. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL_PREFIX__.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>myhouse.cli</string>
    <string>serve</string>
    <string>--config</string>
    <string>__CONFIG__</string>
    <string>--port</string>
    <string>__PORT__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT__</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>Asia/Seoul</string>
    <key>PATH</key>
    <string>__VENV_BIN__:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/dashboard.out</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/dashboard.err</string>
</dict>
</plist>
